# config-vm/qubes-misc-post.service.d/30_qubes-kvm.conf
[Unit]
# XXX: Remove 'qubes-network.service' and 'qubes-firewall.service' depend as
#      they are disabled until they work with KVM.
After=
After=network-pre.target qubes-mount-dirs.service
